# claude-code/rootfs/etc/s6-overlay/s6-rc.d/claude-code/finish
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: Claude Code
# Take down web terminal
# ==============================================================================

bashio::log.info "Claude Code web terminal stopped"
