# openclaw/rootfs/etc/s6-overlay/s6-rc.d/openclaw-init/run
#!/command/with-contenv bashio
# ==============================================================================
# OpenClaw Init Script
# Initializes OpenClaw configuration and installs Chromium
# ==============================================================================

bashio::log.info "Initializing OpenClaw..."

# Create persistent directories
mkdir -p /data/openclaw-config
mkdir -p /data/openclaw-workspace
mkdir -p /data/openclaw-home

# OpenClaw expects config at ~/.openclaw - symlink to persistent storage
# Remove existing directory/symlink if it exists
rm -rf /home/node/.openclaw
ln -sf /data/openclaw-config /home/node/.openclaw

# Create workspace symlink
mkdir -p /data/openclaw-config/workspace
ln -sf /data/openclaw-workspace /data/openclaw-config/workspace

chown -R node:node /data/openclaw-*
chown -h node:node /home/node/.openclaw

# Generate gateway token if not provided
GATEWAY_TOKEN=$(bashio::config 'gateway_token')
if [ -z "$GATEWAY_TOKEN" ]; then
    bashio::log.info "Generating gateway token..."
    GATEWAY_TOKEN=$(openssl rand -hex 32)
    echo "$GATEWAY_TOKEN" > /data/openclaw-config/gateway.token
    bashio::log.info "Gateway token generated and saved"
else
    echo "$GATEWAY_TOKEN" > /data/openclaw-config/gateway.token
fi

# Read configuration
TELEGRAM_BOT_TOKEN=$(bashio::config 'telegram_bot_token')
CLAUDE_SESSION_KEY=$(bashio::config 'claude_session_key')
TELEGRAM_USER_ID=$(bashio::config 'telegram_user_id')
LOG_LEVEL=$(bashio::config 'log_level')

# Create OpenClaw configuration using a heredoc for clean JSON
bashio::log.info "Creating OpenClaw configuration..."

# Build config JSON
cat > /data/openclaw-config/config.json <<JSONEOF
{
  "gateway": {
    "mode": "local",
    "port": 18800,
    "bind": "lan",
    "token": "${GATEWAY_TOKEN}"
  },
  "channels": {
JSONEOF

# Add Telegram block if bot token is provided
if [ -n "$TELEGRAM_BOT_TOKEN" ]; then
    if [ -n "$TELEGRAM_USER_ID" ]; then
        cat >> /data/openclaw-config/config.json <<JSONEOF
    "telegram": {
      "enabled": true,
      "botToken": "${TELEGRAM_BOT_TOKEN}",
      "dmPolicy": "pairing",
      "allowFrom": ["${TELEGRAM_USER_ID}"]
    }
JSONEOF
    else
        cat >> /data/openclaw-config/config.json <<JSONEOF
    "telegram": {
      "enabled": true,
      "botToken": "${TELEGRAM_BOT_TOKEN}",
      "dmPolicy": "pairing"
    }
JSONEOF
    fi
fi

cat >> /data/openclaw-config/config.json <<JSONEOF
  },
  "agents": {
    "defaults": {
      "provider": "anthropic",
      "model": "claude-opus-4-5"
    }
  },
  "logging": {
    "level": "${LOG_LEVEL}"
  }
}
JSONEOF

chown node:node /data/openclaw-config/config.json
chown node:node /data/openclaw-config/gateway.token

bashio::log.info "Configuration written to /home/node/.openclaw/config.json"

# Log the config for debugging (mask sensitive values)
bashio::log.info "Gateway mode: local"
bashio::log.info "Telegram enabled: $([ -n "$TELEGRAM_BOT_TOKEN" ] && echo 'true' || echo 'false')"
bashio::log.info "Claude session: $([ -n "$CLAUDE_SESSION_KEY" ] && echo 'configured' || echo 'not set')"

# Run doctor --fix to auto-enable configured plugins (Telegram, etc.)
bashio::log.info "Running openclaw doctor --fix..."
cd /app && su - node -s /bin/bash -c "cd /app && HOME=/home/node node dist/index.js doctor --fix" || {
    bashio::log.warning "openclaw doctor --fix had warnings, continuing..."
}

# Install Chromium for Playwright (power-user feature)
if [ ! -d "/data/openclaw-home/.cache/ms-playwright" ]; then
    bashio::log.info "Installing Chromium for Playwright (first run)..."
    su - node -c "cd /app && PLAYWRIGHT_BROWSERS_PATH=/data/openclaw-home/.cache/ms-playwright node node_modules/playwright-core/cli.js install chromium" || {
        bashio::log.warning "Playwright Chromium installation failed, but continuing..."
    }
else
    bashio::log.info "Playwright Chromium already installed, skipping..."
fi

# Install additional packages if specified
if bashio::config.has_value 'additional_packages'; then
    bashio::log.info "Installing additional packages..."
    PACKAGES=$(bashio::config 'additional_packages | join(" ")')
    if [ -n "$PACKAGES" ]; then
        apt-get update
        # shellcheck disable=SC2086
        apt-get install -y --no-install-recommends $PACKAGES
        apt-get clean
        rm -rf /tmp/* /var/cache/* /var/lib/apt/lists/*
    fi
fi

bashio::log.info "OpenClaw initialization complete!"
