# openclaw-lisa/rootfs/etc/s6-overlay/s6-rc.d/openclaw-nginx/run
#!/command/with-contenv bashio
# ==============================================================================
# OpenClaw Nginx Proxy Service
# Reverse proxy for WebSocket support through HA ingress
# ==============================================================================

bashio::log.info "Starting nginx reverse proxy..."

exec nginx -c /etc/nginx/nginx.conf
