# openclaw/rootfs/etc/s6-overlay/s6-rc.d/init-openclaw-persist/run
#!/command/with-contenv bashio
# ==============================================================================
# Persist all user config across reboots
# Symlinks ~/.config and ~/.local to /data/openclaw-home so any CLI
# that stores config there (gh, railway, docker, gcloud, etc.) persists.
# ==============================================================================

bashio::log.info "Setting up persistent config symlinks..."

# Ensure persistent dirs exist
mkdir -p /data/openclaw-home/.config /data/openclaw-home/.local/bin

# Replace ephemeral dirs with symlinks to persistent storage
rm -rf /home/node/.config /home/node/.local
ln -sf /data/openclaw-home/.config /home/node/.config
ln -sf /data/openclaw-home/.local /home/node/.local

# Fix ownership
chown -R node:node /data/openclaw-home

bashio::log.info "Persistent config symlinks ready (~/.config, ~/.local)"
