# openclaw/rootfs/etc/s6-overlay/s6-rc.d/openclaw-gateway/run
#!/command/with-contenv bashio
# ==============================================================================
# OpenClaw Gateway Service
# Starts the OpenClaw gateway
# ==============================================================================

bashio::log.info "Starting OpenClaw Gateway..."

# Export environment variables
export HOME=/home/node
export TERM=xterm-256color
export NODE_ENV=production
export PATH="/home/node/.bun/bin:/home/node/.local/bin:${PATH}"

# Gateway configuration
export OPENCLAW_GATEWAY_TOKEN=$(cat /data/openclaw-config/gateway.token)
export OPENCLAW_GATEWAY_BIND="lan"
export OPENCLAW_GATEWAY_PORT="18800"
export OPENCLAW_BRIDGE_PORT="18790"

# Telegram configuration
if bashio::config.has_value 'telegram_bot_token'; then
    export TELEGRAM_BOT_TOKEN=$(bashio::config 'telegram_bot_token')
    bashio::log.info "Telegram bot configured"
fi

# Claude/Anthropic API authentication
# Supports either a real API key (sk-ant-api03-...) or a Claude Pro OAuth token (sk-ant-oat01-...)
if bashio::config.has_value 'claude_session_key'; then
    CLAUDE_SESSION_KEY=$(bashio::config 'claude_session_key')
    case "$CLAUDE_SESSION_KEY" in
        sk-ant-oat01-*)
            # Claude Pro OAuth token — must use CLAUDE_CODE_OAUTH_TOKEN, not ANTHROPIC_API_KEY
            export CLAUDE_CODE_OAUTH_TOKEN="$CLAUDE_SESSION_KEY"
            bashio::log.info "Claude Pro OAuth token configured (CLAUDE_CODE_OAUTH_TOKEN)"
            ;;
        *)
            # Regular Anthropic API key
            export ANTHROPIC_API_KEY="$CLAUDE_SESSION_KEY"
            bashio::log.info "Anthropic API key configured"
            ;;
    esac
fi

# Claude Code OAuth token (for claude CLI coding agent)
# Only set from file if not already configured via claude_session_key above
if [ -z "$CLAUDE_CODE_OAUTH_TOKEN" ] && [ -f /data/openclaw-config/credentials/claude-code-token.txt ]; then
    export CLAUDE_CODE_OAUTH_TOKEN=$(cat /data/openclaw-config/credentials/claude-code-token.txt)
    bashio::log.info "Claude Code OAuth token configured (from credentials file)"
fi

# OpenAI API key (for GPT models and Whisper transcription)
if bashio::config.has_value 'openai_api_key'; then
    export OPENAI_API_KEY=$(bashio::config 'openai_api_key')
    bashio::log.info "OpenAI API key configured"
fi

# Power-user features (enabled by default)
export OPENCLAW_CONFIG_DIR="/data/openclaw-config"
export OPENCLAW_WORKSPACE_DIR="/data/openclaw-workspace"

# Explicit plugin directory — bypasses isSourceCheckoutRoot detection and
# dist-runtime path, using compiled dist/extensions directly (upstream best practice)
export OPENCLAW_BUNDLED_PLUGINS_DIR="/app/dist/extensions"

# Playwright configuration
export PLAYWRIGHT_BROWSERS_PATH="/data/openclaw-home/.cache/ms-playwright"

# Fix ownership on config files modified by root (TTS migration above)
# Avoid chown -R on large workspace dirs — init already handles that
chown node:node /data/openclaw-config/openclaw.json 2>/dev/null || true

# Start the gateway as node user
bashio::log.info "Starting gateway on internal port ${OPENCLAW_GATEWAY_PORT} (nginx proxy on 18789)..."

cd /app || exit 1

# ---- Migrate legacy TTS provider keys (runs every gateway start) ----
# Strict schema rejects openai/edge/elevenlabs/microsoft as direct children
# of messages.tts or channels.discord.voice.tts — must be under .providers.
CONFIG="/data/openclaw-config/openclaw.json"
if [ -f "$CONFIG" ]; then
    jq '
      (if .messages.tts and ([.messages.tts | has("openai"), has("edge"), has("elevenlabs"), has("microsoft")] | any) then
        .messages.tts |= (
          (if has("openai") then .providers.openai = (.providers.openai // {}) * .openai | del(.openai) else . end) |
          (if has("edge") then .providers.microsoft = (.providers.microsoft // {}) * .edge | del(.edge) else . end) |
          (if has("elevenlabs") then .providers.elevenlabs = (.providers.elevenlabs // {}) * .elevenlabs | del(.elevenlabs) else . end) |
          (if has("microsoft") then .providers.microsoft = (.providers.microsoft // {}) * .microsoft | del(.microsoft) else . end)
        )
      else . end) |
      (if .channels.discord.voice.tts and ([.channels.discord.voice.tts | has("openai"), has("edge"), has("elevenlabs"), has("microsoft")] | any) then
        .channels.discord.voice.tts |= (
          (if has("openai") then .providers.openai = (.providers.openai // {}) * .openai | del(.openai) else . end) |
          (if has("edge") then .providers.microsoft = (.providers.microsoft // {}) * .edge | del(.edge) else . end) |
          (if has("elevenlabs") then .providers.elevenlabs = (.providers.elevenlabs // {}) * .elevenlabs | del(.elevenlabs) else . end) |
          (if has("microsoft") then .providers.microsoft = (.providers.microsoft // {}) * .microsoft | del(.microsoft) else . end)
        )
      else . end)
    ' "$CONFIG" > /tmp/oc-tts-gw.json && mv /tmp/oc-tts-gw.json "$CONFIG" && chown node:node "$CONFIG" \
        && bashio::log.info "TTS config migration applied" \
        || bashio::log.warning "TTS config migration skipped (jq error)"
fi

# Run gateway with error capture — log the crash reason and exit code
# (s6 will restart this longrun service on exit)
s6-setuidgid node \
    node dist/index.js gateway \
    --allow-unconfigured \
    --bind "${OPENCLAW_GATEWAY_BIND}" \
    --port "${OPENCLAW_GATEWAY_PORT}" 2>&1

EXIT_CODE=$?
bashio::log.error "Gateway exited with code ${EXIT_CODE}"

# Brief pause to prevent rapid crash-restart cycling and allow log to flush
sleep 10
