# channels-dvr/rootfs/etc/services.d/channels-dvr/run
#!/command/with-contenv bashio

bashio::log.info "Starting Channels DVR Server..."

# Define paths
SEED_DIR="/opt/channels-dvr"
BINARY_DIR="/media/dvr/channels-dvr-server"
DATA_DIR="/media/dvr/channels-dvr-data"

# Ensure binary directory is initialized
if [ ! -d "$BINARY_DIR/latest" ]; then
    bashio::log.info "First run detected. Initializing binary storage at $BINARY_DIR..."
    # Only create directory if parent mount exists and this directory doesn't
    if [ -d "/media/dvr" ] && [ ! -d "$BINARY_DIR" ]; then
        mkdir -p "$BINARY_DIR"
    fi
    cp -a "$SEED_DIR/." "$BINARY_DIR/"
    bashio::log.info "Binary initialization complete."
else
    bashio::log.info "Found existing binary at $BINARY_DIR"
fi

# Ensure data directory exists (only create if NFS mount is present but dir is missing)
if [ -d "/media/dvr" ] && [ ! -d "$DATA_DIR" ]; then
    mkdir -p "$DATA_DIR"
fi
bashio::log.info "Using data directory: $DATA_DIR"

# Debug: show what's in the data directory
bashio::log.info "DEBUG: Contents of $DATA_DIR:"
ls -la "$DATA_DIR/" 2>&1 | while read line; do bashio::log.info "  $line"; done
bashio::log.info "DEBUG: Mount info:"
mount | grep -E "media|dvr|nfs" 2>&1 | while read line; do bashio::log.info "  $line"; done
bashio::log.info "DEBUG: df for DATA_DIR:"
df -h "$DATA_DIR" 2>&1 | while read line; do bashio::log.info "  $line"; done

# Auto-restore database from backup if missing or empty
# The backup settings.db is ~60KB; a fresh/empty one is much smaller
SETTINGS_SIZE=0
if [ -f "$DATA_DIR/settings.db" ]; then
    SETTINGS_SIZE=$(stat -c%s "$DATA_DIR/settings.db" 2>/dev/null || echo 0)
fi
if [ "$SETTINGS_SIZE" -lt 50000 ]; then
    bashio::log.warning "Database missing or empty in $DATA_DIR (size: ${SETTINGS_SIZE}B)! Attempting restore from backup..."
    # Try NAS backup first (has the good 204KB recorder.db)
    NAS_BACKUP_BASE="/media/dvr/Database"
    BACKUP_BASE="/share/channels-dvr/Database"
    if [ -d "$NAS_BACKUP_BASE" ]; then
        LATEST_BACKUP=$(ls -td "$NAS_BACKUP_BASE"/backup-*/ 2>/dev/null | head -1)
        if [ -n "$LATEST_BACKUP" ] && [ -f "$LATEST_BACKUP/settings.db" ]; then
            bashio::log.info "Restoring from NAS backup: $LATEST_BACKUP"
            cp "$LATEST_BACKUP/settings.db" "$DATA_DIR/"
            [ -f "$LATEST_BACKUP/recorder.db" ] && cp "$LATEST_BACKUP/recorder.db" "$DATA_DIR/"
            bashio::log.info "Database restored from NAS successfully."
            SETTINGS_SIZE=$(stat -c%s "$DATA_DIR/settings.db" 2>/dev/null || echo 0)
        fi
    fi
    # Fall back to local backup if NAS restore didn't work
    if [ "$SETTINGS_SIZE" -lt 50000 ]; then
        bashio::log.warning "NAS backup not available, trying local backup..."
        BACKUP_BASE="/share/channels-dvr/Database"
    if [ -d "$BACKUP_BASE" ]; then
        LATEST_BACKUP=$(ls -td "$BACKUP_BASE"/backup-*/ 2>/dev/null | head -1)
        if [ -n "$LATEST_BACKUP" ] && [ -f "$LATEST_BACKUP/settings.db" ]; then
            bashio::log.info "Restoring from backup: $LATEST_BACKUP"
            cp "$LATEST_BACKUP/settings.db" "$DATA_DIR/"
            [ -f "$LATEST_BACKUP/recorder.db" ] && cp "$LATEST_BACKUP/recorder.db" "$DATA_DIR/"
            bashio::log.info "Database restored successfully."
        else
            bashio::log.warning "No valid backup found in $BACKUP_BASE"
        fi
    else
        bashio::log.warning "No backup directory found at $BACKUP_BASE"
    fi
    fi
else
    bashio::log.info "Database found in $DATA_DIR (size: ${SETTINGS_SIZE}B)"
fi

# Create default recordings directory if it doesn't exist
if [ -d "/share" ]; then
    mkdir -p "/share/channels-dvr"
fi

# Note: /media/dvr is an NFS mount - do NOT create it here
# It should be mounted by the host system before addon starts

# Start Channels DVR Server with explicit data path
cd "$BINARY_DIR/latest"

bashio::log.info "Starting Channels DVR binary from $(pwd) with data path $DATA_DIR"
exec ./channels-dvr -dir "$DATA_DIR"
