# channels-dvr/rootfs/etc/services.d/channels-dvr/run
#!/command/with-contenv bashio

bashio::log.info "Starting Channels DVR Server..."

# Define paths
SEED_DIR="/opt/channels-dvr"
BINARY_DIR="/media/dvr/channels-dvr-server"
DATA_DIR="/media/dvr/channels-dvr-data"

# Ensure binary directory is initialized
if [ ! -d "$BINARY_DIR/latest" ]; then
    bashio::log.info "First run detected. Initializing binary storage at $BINARY_DIR..."
    mkdir -p "$BINARY_DIR"
    cp -a "$SEED_DIR/." "$BINARY_DIR/"
    bashio::log.info "Binary initialization complete."
else
    bashio::log.info "Found existing binary at $BINARY_DIR"
fi

# Ensure data directory exists
mkdir -p "$DATA_DIR"
bashio::log.info "Using data directory: $DATA_DIR"

# Create default recordings directory if it doesn't exist
if [ -d "/share" ]; then
    mkdir -p "/share/channels-dvr"
fi

if [ -d "/media" ]; then
    mkdir -p "/media/dvr"
fi

# Start Channels DVR Server with explicit data path
cd "$BINARY_DIR/latest"

bashio::log.info "Starting Channels DVR binary from $(pwd) with data path $DATA_DIR"
exec ./channels-dvr -path "$DATA_DIR"
