# claude-code/rootfs/etc/s6-overlay/s6-rc.d/claude-remote/finish
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: Claude Code
# Take down claude remote-control
# ==============================================================================

# Kill any lingering remote-control processes
pkill -f "claude remote-control" 2>/dev/null || true

bashio::log.info "Claude remote-control stopped"
