# claude-code/rootfs/etc/s6-overlay/s6-rc.d/claude-remote/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: Claude Code
# Runs claude remote-control for one or more workspace directories
# ==============================================================================

if ! bashio::config.true 'auto_start_remote_control'; then
    bashio::log.info "Claude remote-control auto-start is disabled"
    exec sleep infinity
fi

sleep 3

AUTH_MODE=$(bashio::config 'auth_mode')

# For subscription mode, wait until the user has logged in
if [[ "${AUTH_MODE}" == "subscription" ]]; then
    if [[ ! -f "/data/claude-code/.claude/credentials.json" ]] && \
       [[ ! -f "/data/claude-code/.claude/.credentials.json" ]] && \
       [[ ! -d "/data/claude-code/.claude/auth" ]]; then
        bashio::log.warning "Waiting for subscription login before starting remote-control..."
        bashio::log.warning "Open the web terminal and run: claude login"
        while [[ ! -f "/data/claude-code/.claude/credentials.json" ]] && \
              [[ ! -f "/data/claude-code/.claude/.credentials.json" ]] && \
              [[ ! -d "/data/claude-code/.claude/auth" ]]; do
            sleep 10
        done
        bashio::log.info "Login detected! Starting remote-control..."
        sleep 2
    fi
fi

if [[ "${AUTH_MODE}" == "api_key" ]] && ! bashio::config.has_value 'anthropic_api_key'; then
    bashio::log.error "No API key configured - remote-control cannot start"
    exec sleep infinity
fi

# --- Build the list of workspaces ---
declare -a WORKSPACES=()

if bashio::config.has_value 'remote_workspaces'; then
    REMOTE_WS_COUNT=$(bashio::config 'remote_workspaces | length')
    for (( i=0; i<REMOTE_WS_COUNT; i++ )); do
        WS=$(bashio::config "remote_workspaces[${i}]")
        if [[ -d "${WS}" ]]; then
            WORKSPACES+=("${WS}")
        else
            bashio::log.warning "Skipping remote workspace '${WS}' - directory not found"
        fi
    done
fi

# Fallback: if no remote_workspaces configured (or none valid), use single workspace
if [[ ${#WORKSPACES[@]} -eq 0 ]]; then
    WORKSPACE=$(bashio::config 'workspace')
    WORKSPACE="${WORKSPACE:-/share/code}"
    if [[ -d "${WORKSPACE}" ]]; then
        WORKSPACES+=("${WORKSPACE}")
    else
        WORKSPACES+=("/homeassistant")
    fi
fi

# --- Read permission mode ---
PERMISSION_MODE="bypassPermissions"
if bashio::config.has_value 'permission_mode'; then
    PERMISSION_MODE=$(bashio::config 'permission_mode')
fi

bashio::log.info "Starting Claude remote-control for ${#WORKSPACES[@]} workspace(s) (mode: ${PERMISSION_MODE})"

# --- Per-workspace process tracking ---
declare -A PID_TO_WS=()
declare -A PID_TO_NAME=()

start_remote() {
    local ws_path="$1"
    local ws_name
    ws_name=$(basename "${ws_path}")

    bashio::log.info "Starting remote-control: ${ws_name} (${ws_path})"
    (
        cd "${ws_path}" || exit 1
        exec sh -c "echo y | claude remote-control --name '${ws_name}' --permission-mode '${PERMISSION_MODE}' > /dev/null"
    ) &
    local pid=$!
    PID_TO_WS[${pid}]="${ws_path}"
    PID_TO_NAME[${pid}]="${ws_name}"
    bashio::log.debug "  PID ${pid} -> ${ws_name}"
}

# --- Cleanup on signal ---
cleanup() {
    bashio::log.info "Stopping all remote-control processes..."
    for pid in "${!PID_TO_WS[@]}"; do
        kill "${pid}" 2>/dev/null
    done
    wait
    exit 0
}
trap cleanup SIGTERM SIGINT

# --- Start all workspaces ---
for ws in "${WORKSPACES[@]}"; do
    start_remote "${ws}"
done

# --- Supervision loop: restart any child that exits ---
while true; do
    wait -n 2>/dev/null
    exit_code=$?

    # Find which PID(s) are no longer running and restart them
    for pid in "${!PID_TO_WS[@]}"; do
        if ! kill -0 "${pid}" 2>/dev/null; then
            ws_path="${PID_TO_WS[${pid}]}"
            ws_name="${PID_TO_NAME[${pid}]}"
            unset "PID_TO_WS[${pid}]"
            unset "PID_TO_NAME[${pid}]"
            bashio::log.warning "Remote-control '${ws_name}' exited (code ${exit_code}), restarting in 5s..."
            sleep 5
            if [[ -d "${ws_path}" ]]; then
                start_remote "${ws_path}"
            else
                bashio::log.error "Workspace '${ws_path}' no longer exists, not restarting"
            fi
        fi
    done

    # If no children remain, exit (s6 will restart the whole service)
    if [[ ${#PID_TO_WS[@]} -eq 0 ]]; then
        bashio::log.error "No remote-control processes running, exiting"
        exit 1
    fi
done
